# launchd/com.podcast-updates.finance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Finance pipeline: Saturday 9am PT.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.podcast-updates.finance</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/USERNAME/Podcast-Updates/.venv/bin/podcast-updates</string>
    <string>--config</string>
    <string>/Users/USERNAME/Podcast-Updates/config/shows_finance.yaml</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/USERNAME/Podcast-Updates</string>

  <!-- Saturday at 9:00 AM local time -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/USERNAME/Podcast-Updates/data/logs/finance.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/USERNAME/Podcast-Updates/data/logs/finance.log</string>
</dict>
</plist>
