# launchd/com.podcast-updates.tech.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Tech pipeline: Friday 5pm PT.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.podcast-updates.tech</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/USERNAME/Podcast-Updates/.venv/bin/podcast-updates</string>
    <string>--config</string>
    <string>/Users/USERNAME/Podcast-Updates/config/shows_tech.yaml</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/USERNAME/Podcast-Updates</string>

  <!-- Friday at 5:00 PM local time -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/USERNAME/Podcast-Updates/data/logs/tech.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/USERNAME/Podcast-Updates/data/logs/tech.log</string>
</dict>
</plist>
